# ansible/playbook/roles/nsq/templates/etc/init/nsqadmin.conf
description "nsqadmin"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]
limit nofile 50000 50000

respawn
respawn limit 10 5

pre-start script
    test -x {{ nsq_install_directory }}/{{ nsq_release }}/bin/nsqadmin || { stop; exit 0; }
end script

exec {{ nsq_install_directory }}/{{ nsq_release }}/bin/nsqd {{ nsq_nsqd_opts }} {% if nsqlookupd_role and nsq_lookupd_tcp_addresses | length > 0 %}{% for item in nsq_lookupd_tcp_addresses %}-lookupd-tcp-address={{ item }}:{{ nsq_lookupd_tcp_port }} {% endfor %}{% endif %}
